# conf/linkace.cron
*   *   *   *   *   cd __INSTALL_DIR__ && php__PHP_VERSION__ artisan schedule:run > /dev/null